# src/systemd/units/arkloned-quake.path
[Path]
PathChanged=/roms/ports/quake2/baseq2/save
Unit=arkloned@-roms-ports-quake2-baseq2-save\x40quake2-baseq2-save\x40.service

[Install]
WantedBy=multi-user.target
