# src/systemd/units/arkloned-amiberry.path
[Path]
PathChanged=/opt/amiberry/savestates
Unit=arkloned@-opt-amiberry-savestates\x40amiberry-savestates\x40.service

[Install]
WantedBy=multi-user.target
