# systemd/units/arkloned-ppsspp.path
[Path]
PathChanged=/home/ark/.config/ppsspp
Unit=arkloned@-home-ark-.config-ppsspp\x40ppsspp\x40ppsspp.service

[Install]
WantedBy=multi-user.target
