# systemd/units/arkloned-cavestory.path
[Path]
PathChanged=/roms/ports/cavestory
Unit=arkloned@-roms-ports-cavestory\x40cavestory\x40cavestory.service

[Install]
WantedBy=multi-user.target
